# systemd/neonsignal-redirect.service
[Unit]
Description=neonsignal redirector (80 -> 443)
After=network-online.target
Wants=network-online.target

[Service]
User=core
Group=core
WorkingDirectory=/home/core/code/neonsignal
Environment=REDIRECT_INSTANCES=3
Environment=REDIRECT_PORT=9090
Environment=REDIRECT_TARGET_PORT=443
Environment=REDIRECT_HOST=10.0.0.10
ExecStart=/usr/local/bin/neonsignal_redirect
Restart=on-failure
AmbientCapabilities=CAP_NET_BIND_SERVICE
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
